# flask-app/endpoint/config/supervisor.conf
[program:api-uwsgi]
command = /usr/local/bin/uwsgi --http :5000 --wsgi-file /srv/api/endpoint.py

