# docker-storm/supervisor.conf
[program:start-topology]
environment LEIN_ROOT="yes"
directory = /srv/toycount-topology
command = /usr/local/bin/sparse run 

